# temp/dashboard-rotator-kiosk.service
[Unit]
Description=Dashboard Rotator Kiosk Display
After=dashboard-rotator-server.service graphical.target
Wants=dashboard-rotator-server.service

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
ExecStartPre=/usr/bin/sleep 5
ExecStart=/home/pi/dashboard-rotator/start-kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target
